# challenges/attack-defense/src/supervisord.conf
# supervisord runs as PID 1 so the service runs as a CHILD: if it crashes (e.g. a
# botched exploit kills it) supervisord auto-restarts it in place — autorestart +
# effectively infinite retries — WITHOUT the container exiting, so the box stays up
# and keeps the team's filesystem changes (patches). NOTE: the platform exposes no
# in-place restart; its operator "Force Restart" / team "Reset" destroys and
# recreates the container from this image (new IP), losing all patches.
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:svc]
command=python3 /service.py
autostart=true
autorestart=true
startretries=1000000
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
